# deploy/supervisord.conf
; Supervisor configuration for non-Docker deployments (architecture doc — Deployment & DevOps).
; Assumes the repo lives at /app with frontend built to /app/frontend/dist
; and Python deps installed for the backend.

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:backend]
command=uvicorn server:app --host 0.0.0.0 --port 8000
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/backend.err.log
stdout_logfile=/var/log/backend.out.log

; Serve the built frontend with any static server (example: `npx serve`).
; In production prefer the nginx setup in frontend/deploy/nginx.conf.
[program:frontend]
command=npx serve -s /app/frontend/dist -l 3000
directory=/app/frontend
autostart=true
autorestart=true
stderr_logfile=/var/log/frontend.err.log
stdout_logfile=/var/log/frontend.out.log
